# scripts/com.market-warehouse.daily-update-watchdog.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.market-warehouse.daily-update-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <!-- Replace /path/to/repo with your actual repo path -->
        <string>/path/to/repo/scripts/run_daily_update_watchdog.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- 18:30 PT leaves room for the 13:05 PT sync, retries, and a long-running close-of-day ingest -->
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/com.market-warehouse.daily-update-watchdog.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.market-warehouse.daily-update-watchdog.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
