# scripts/com.market-warehouse.daily-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.market-warehouse.daily-update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>$HOME/dev/apps/finance/market-data-warehouse/scripts/run_daily_update.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- 21:05 UTC = 4:05 PM EST / 5:05 PM EDT — always after market close -->
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/com.market-warehouse.daily-update.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.market-warehouse.daily-update.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
